# etc/init/fontello.conf
# fontello.com upstart script
description "Fontello - Iconic fonts scissors"
author "Nodeca"

start on startup
stop on shutdown

env NODECA_ENV="production"

script
    cd /path/to/your/script/
    exec node fontomas.js 2>&1
end script
